# etc/sidewinderd.service.in
[Unit]
Description=Support for Microsoft SideWinder X4 / X6 and Logitech G105 / G710+
After=multi-user.target

[Service]
ExecStart=${CMAKE_INSTALL_PREFIX}/bin/sidewinderd

[Install]
WantedBy=multi-user.target
